# fstab/mnt-disk4.automount
[Unit]
Description=Automount mergerfs disk4

[Automount]
Where=/mnt/disk4
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target